# firewall/etc/network/if-up.d/iptables
#!/bin/sh
/etc/firewall.conf
/etc/getMasterFirewall.php&